# high-level/rosbuild/scout_gumstix/etc/init/roscore.conf
start on runlevel [12345]
stop on runlevel [!12345]

setuid scout
setgid scout

respawn
console log

script
  . /opt/ros/fuerte/setup.sh
  exec roscore
end script
